# backend/docker/back4app-supervisord.conf
[program:web]
command=apache2-foreground
priority=10
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue]
command=/bin/sh -c 'if [ "$BACK4APP_RUN_QUEUE_WORKER" = "true" ]; then exec gosu www-data php artisan queue:work --queue=default --sleep=3 --tries=3 --backoff=5 --timeout=90 --max-time=3600; else exec sleep infinity; fi'
directory=/var/www/html
priority=20
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=BACK4APP_RUN_QUEUE_WORKER="%(ENV_BACK4APP_RUN_QUEUE_WORKER)s"

[program:scheduler]
command=/bin/sh -c 'if [ "$BACK4APP_RUN_SCHEDULER" = "true" ]; then exec gosu www-data php artisan schedule:work; else exec sleep infinity; fi'
directory=/var/www/html
priority=30
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=90
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=BACK4APP_RUN_SCHEDULER="%(ENV_BACK4APP_RUN_SCHEDULER)s"
